# fs/etc/services.d/dhcp-exporter/run
#!/usr/bin/with-contenv sh

exec runas /usr/bin/dhcp-exporter $DHCP_EXPORTER_ARGS